# device-clamav-el9/src/main/resources/etc/systemd/system/clamd@scan.service.d/local.conf
[Unit]
PartOf=clamav-autodiscovery.service

